# supervisor/vagrant.conf
[program:vagrant]
directory=/opt/vagrant/src
command=python vagrant.py 80
